# skill/folio-scribe/scripts/launchd/com.folio-scribe.weekly-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.folio-scribe.weekly-summary</string>
    <key>ProgramArguments</key>
    <array>
        <string>__SKILL_DIR__/scripts/write_periodic_summary.py</string>
        <string>--vault</string>
        <string>__VAULT__</string>
        <string>--period</string>
        <string>weekly</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>__WEEKLY_SUMMARY_HOUR__</integer>
        <key>Minute</key>
        <integer>__WEEKLY_SUMMARY_MINUTE__</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/launchd-weekly-summary.out</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/launchd-weekly-summary.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Library/Frameworks/Python.framework/Versions/3.13/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>FOLIO_SCRIBE_AI_CLI</key>
        <string>__AI_CLI__</string>
        <key>FOLIO_SCRIBE_AI_FALLBACK_CLI</key>
        <string>__AI_FALLBACK_CLI__</string>
        <key>FOLIO_SCRIBE_AI_MODEL</key>
        <string>__AI_MODEL__</string>
        <key>FOLIO_SCRIBE_WEB_EXPORT_DIR</key>
        <string>__WEB_EXPORT_DIR__</string>
        <key>FOLIO_SCRIBE_WEB_TITLE</key>
        <string>__WEB_TITLE__</string>
        <key>FOLIO_SCRIBE_WEB_DEPLOY</key>
        <string>__WEB_DEPLOY__</string>
        <key>FOLIO_SCRIBE_VERCEL_PROJECT</key>
        <string>__VERCEL_PROJECT__</string>
        <key>FOLIO_SCRIBE_VERCEL_SCOPE</key>
        <string>__VERCEL_SCOPE__</string>
        <key>FOLIO_SCRIBE_VERCEL_TARGET</key>
        <string>__VERCEL_TARGET__</string>
        <key>FOLIO_SCRIBE_VERCEL_DOMAIN</key>
        <string>__VERCEL_DOMAIN__</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>__VAULT__</string>
</dict>
</plist>
